# beedude.socket
[Unit]
Description=Beedude socket

[Socket]
ListenStream=/run/beedude.sock

[Install]
WantedBy=sockets.target